# .config/sv/mpd/run
#!/bin/sh
exec mpd --no-daemon --stderr ${OPTS:-} 2>&1